# ansible/templates/systemd/ttyd-shell.target
[Unit]
Description=ttyd-shell target

[Install]
WantedBy=default.target
